# app/crontab
# Echo current time and execute /var/www/html/run_cron.
* * * * * root /usr/local/bin/crontask.sh >> /var/log/cron.log 2>&1
